# systemd/riot-broker.service
[Unit]
Description=Riot Broker Application
After=network.target radvd.service

[Service]
User=pi
ExecStart=/usr/bin/python3 /home/pi/demos/riot-broker.py

[Install]
WantedBy=multi-user.target
Alias=riot-broker.service
